# etc/service/iperf/run
#!/bin/bash -e

exec /usr/bin/iperf -s -p 5001
